# net-misc/openlldp/files/openlldp.initd
#!/sbin/runscript

depend() {
	need net
}


start() {
	ebegin "Starting LLDP"
	start-stop-daemon --start --quiet --background  --exec /usr/lib/openlldp/sbin/lldpd -- ${OPTS} >/dev/null
	eend $?
}

stop() {
	ebegin "Stopping LLDP"
	start-stop-daemon --stop --quiet --retry 20 --exec /usr/lib/openlldp/sbin/lldpd
	eend $?
}
