# net-wireless/bluez-utils/files/bluez-utils-2.24-udev.rules
# Bluetooth devices:
#
KERNEL=="hci[0-9]*", RUN+="/lib/udev/bluetooth.sh"
