# net-wireless/bluez-utils/files/bluez-utils-3.10-init.d-sdpd
#!/sbin/runscript
# Copyright 1999-2007 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

depend() {
	need hcid
}

start() {
	ebegin "Starting sdpd"
	start-stop-daemon --start --quiet --exec /usr/sbin/sdpd
	eend $?
}

stop() {
	ebegin "Stopping sdpd"
	start-stop-daemon --stop --quiet --exec /usr/sbin/sdpd
	eend $?
}

