# net-p2p/freenet/files/freenet
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

start() {
	if ! [[ -e /opt/freenet/freenet.ini ]]; then
		start-stop-daemon --start --chuid freenet --exec /bin/sh /var/freenet/run.sh start 2>/dev/null
		eend 0
	else
		echo "Please move freenet to the new location with the following command:"
		echo "         mv /opt/freenet /var/freenet"
		eend 1
	fi
}

status() {
	start-stop-daemon --start --chuid freenet --exec /bin/sh /var/freenet/run.sh status 2>/dev/null
	eend 0
}

stop() {
	start-stop-daemon --start --chuid freenet --exec /bin/sh /var/freenet/run.sh stop 2>/dev/null
	eend 0
}
