# dev-db/membase-server-community/files/1.6.0.1/init-epmd
#!/sbin/runscript
# Copyright 1999-2010 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

VER=1.6.0.1
PIDFILE=/var/run/membase-epmd.pid
MEMBASE_USER=membase
MEMBASE_PATH=/opt/membase/$VER

depend() {
	need net
}

start() {
	ebegin "Starting membase (epmd)"
	cd $MEMBASE_PATH/ns_server

        touch $PIDFILE
	chown membase:daemon $PIDFILE

	# Initialize distributed erlang on the system (i.e. epmd)
	start-stop-daemon --start \
		--oknodo -c membase -m --pidfile $PIDFILE \
		--exec /usr/bin/erl -- \
			-noshell -setcookie nocookie -sname init -run init stop 2>&1 > /dev/null
}

stop() {
	ebegin "Stopping membase (epmd)"
	start-stop-daemon --stop --name epmd
}
