# dev-db/pgbouncer/files/pgbouncer.initd
#!/sbin/runscript

PIDFILE=/var/run/pgbouncer/pgbouncer.pid

opts="start stop reload"

depend() {
	need net
	after postgresql
}

start() {
	ebegin "Starting pgbouncer"
	
	start-stop-daemon --start --pidfile ${PIDFILE} \
		--exec /usr/bin/pgbouncer -- -q -d -u pgbouncer /etc/pgbouncer.conf
	eend $? "Failed to start pgbouncer"
}

stop() {
	ebegin "Stopping pgbouncer"
	start-stop-daemon --stop --quiet --pidfile ${PIDFILE}
	eend $?
} 

reload() {
	ebegin "Gracefully reloading pgbouncer"
	kill -HUP `cat ${PIDFILE}` &>/dev/null
	eend $?
}
