# net-dns/nsd/files/nsd.cron
#!/bin/sh 
# Copyright 1999-2005 Gentoo Foundation 
# Distributed under the terms of the GNU General Public License v2 
# $Header: $ 
# 
NSDC_BIN="/usr/sbin/nsdc" 

if [[ -x ${NSDC_BIN} ]]; then 
	${NSDC_BIN} update >>/dev/null 2>&1 
fi
