# net-dns/nsd/files/nsd.cron
#!/bin/sh 
# 
NSDC_BIN="/usr/sbin/nsdc" 

if [[ -x ${NSDC_BIN} ]]; then 
	${NSDC_BIN} update >>/dev/null 2>&1 
fi
